# systemd/waytype-daemon.service
[Unit]
Description=WayType GUI Daemon
After=graphical-session.target

[Service]
# NOTE: Replace /home/prakhar/Projects/voice-typing with the absolute path to your cloned repository.
ExecStart=/usr/bin/python3 /home/prakhar/Projects/voice-typing/gui_daemon.py
WorkingDirectory=/home/prakhar/Projects/voice-typing
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
